# beta-test/docker/examples/supervisord.conf
; beta-test/docker/examples/supervisord.conf
;
; Qwen3.5-4B (テキスト) + Gemma-4-26B (画像認識) 同時稼働の設定例。
; Example config for running Qwen3.5-4B (text) + Gemma-4-26B (vision) simultaneously.
;
; 使い方 / Usage:
;   このファイルを beta-test/docker/supervisord.conf にコピーして使う。
;   Copy this file to beta-test/docker/supervisord.conf before starting.
;
;   cp beta-test/docker/examples/supervisord.conf beta-test/docker/supervisord.conf
;
; モデルファイルのパス構成 / Model file layout (MODELS_DIR=$HOME/gguf_models):
;
;   ~/gguf_models/
;   └── unsloth/
;       ├── Qwen3.5-4B-GGUF/
;       │   └── Qwen3.5-4B-UD-Q4_K_XL.gguf          ← llm-qwen4b が読む
;       └── gemma-4-26B-A4B-it-GGUF/
;           ├── gemma-4-26B-A4B-it-UD-Q4_K_XL.gguf   ← llm-gemma26b が読む
;           └── mmproj-BF16.gguf                      ← 画像認識用プロジェクション

; ────────────────────────────────────────────────────────────────
; supervisord メイン設定 / Main supervisord settings
; ────────────────────────────────────────────────────────────────
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ────────────────────────────────────────────────────────────────
; llm-qwen4b: Qwen3.5-4B テキストモデル (port 8081)
; VRAM 約 2.8GB
; ────────────────────────────────────────────────────────────────
[program:llm-qwen4b]
command=bash /app/beta-test/docker/scripts/llama-entrypoint.sh
environment=MODEL_FILE="unsloth/Qwen3.5-4B-GGUF/Qwen3.5-4B-UD-Q4_K_XL.gguf",LLAMA_PORT="8081",N_PARALLEL="4",CTX_SIZE="8192",N_GPU_LAYERS="-1",N_PREDICT="512",LD_LIBRARY_PATH="/opt/llama/bin"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/app/logs/llm-qwen4b.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/llm-qwen4b.log
stderr_logfile_maxbytes=0

; ────────────────────────────────────────────────────────────────
; llm-gemma26b: Gemma-4-26B 画像認識モデル (port 8082)
; VRAM 約 13GB / mmproj-BF16.gguf を同時読み込み
; N_PARALLEL を 2 に抑えてVRAMを節約
; ────────────────────────────────────────────────────────────────
[program:llm-gemma26b]
command=bash /app/beta-test/docker/scripts/llama-entrypoint.sh
environment=MODEL_FILE="unsloth/gemma-4-26B-A4B-it-GGUF/gemma-4-26B-A4B-it-UD-Q4_K_XL.gguf",MMPROJ_FILE="unsloth/gemma-4-26B-A4B-it-GGUF/mmproj-BF16.gguf",LLAMA_PORT="8082",N_PARALLEL="2",CTX_SIZE="8192",N_GPU_LAYERS="-1",N_PREDICT="512",LD_LIBRARY_PATH="/opt/llama/bin"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/app/logs/llm-gemma26b.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/llm-gemma26b.log
stderr_logfile_maxbytes=0

; ────────────────────────────────────────────────────────────────
; wrapper: Bun + Hono プロキシサーバー
; ────────────────────────────────────────────────────────────────
[program:wrapper]
command=bun run start
directory=/app
environment=MONGO_URI="%(ENV_MONGO_URI)s",MASTER_KEY="%(ENV_MASTER_KEY)s",BACKENDS_CONFIG="%(ENV_BACKENDS_CONFIG)s",PORT="%(ENV_PORT)s"
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/app/logs/wrapper.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/wrapper.log
stderr_logfile_maxbytes=0
